# appliance/skeleton/etc/init.d/warrior4-appliance-display
#!/sbin/openrc-run

command=/usr/bin/warrior4-appliance-display
command_args="--vt 1"
command_background=true
pidfile=/var/run/warrior4-appliance-display.pid
